# setup/resources/enable-can.service
[Unit]
After=network.target

[Service]
ExecStart=/home/msr/setup/resources/enable-can.sh

[Install]
WantedBy=default.target
